# crontab
# Gap down stocks daily analysis - 9:15 AM ET timestamp + a few minutes
53 02 * * * cd /app && python gap_down_email.py email
50 12 * * * cd /app && python gap_down_email.py email
55 03 * * * cd /app && python gap_down_email.py email
40 04 * * * cd /app && python gap_down_email.py email
51 04 * * * cd /app && python gap_down_email.py email
10 17 * * * cd /app && python gap_down_email.py email
20 13 * * * cd /app && python gap_down_email.py email-all
